# deploy/com.asmi.daemon.launchagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  asmi LaunchAgent — user-domain (gui/$UID)

  Install to ~/Library/LaunchAgents/com.asmi.daemon.plist
  Bootstrap: launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.asmi.daemon.plist

  USE THIS when the node has a graphical login session (console user = ma, not root).
  Check with: stat -f '%Su' /dev/console
  If console shows `root`, use the LaunchDaemon variant instead (needs sudo).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.asmi.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>__ASMI_HOME__/.cargo/bin/asmi</string>
        <string>--serve</string>
        <string>--port</string>
        <string>9090</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__ASMI_HOME__/Library/Logs/asmi-daemon.log</string>
    <key>StandardErrorPath</key>
    <string>__ASMI_HOME__/Library/Logs/asmi-daemon.log</string>
    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
